# packaging/systemd/user/wheelctl.service
[Unit]
Description=wheelctl evdev volume wheel (user service)
Documentation=https://github.com/gracee3/wheelctl-rs
After=pipewire.service wireplumber.service

[Service]
Type=simple
ExecStart=%h/.cargo/bin/wheelctl run
Restart=on-failure
RestartSec=2
Environment=RUST_LOG=wheelctl=info
Nice=10
IOSchedulingClass=idle
CPUWeight=10
IOWeight=10

[Install]
WantedBy=default.target
